# supporting_files/supervisord-mysqld.conf
[program:mysqld]
command=/usr/bin/mysqld_safe --datadir=/var/lib/mysql
numprocs=1
autostart=true
autorestart=true