# supporting_files/supervisord-webmin.conf
[program:webmin]
user=root
command=/usr/sbin/service webmin start
numprocs=1
autostart=true
startsecs=0