# supporting_files/supervisord-webmin.conf
[program:webmin]
user=root
command=/start-webmin.sh
numprocs=1
autostart=true
startsecs=0